# scripts/systemd/mirrormini.service
[Unit]
Description=MirrorMini Startup Service
After=local-fs.target
Before=network.target
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/bin/bash /path/to/your/project/scripts/mirrormini.sh
WorkingDirectory=/path/to/your/project
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

